# contrail/contrail/3.0/latest/analytics-collector/supervisor/collector.conf
[program:contrail-collector]
command=/usr/bin/contrail-collector
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

