# contrail/openstack/kilo/nova-consoleauth/supervisord/openstack.conf
[program:nova-consoleauth]
command=/bin/bash -c "exec start-stop-daemon --start --chuid root --chdir /var/lib/nova --name nova-consoleauth --exec /usr/bin/nova-consoleauth"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
