# contrail/nova-api/supervisord/openstack.conf
[program:nova-api]
command=/bin/bash -c "exec start-stop-daemon --start --chuid root --chdir /var/lib/nova --name nova-api --exec /usr/bin/nova-api"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
